# butterfly.service
[Unit]
Description=Butterfly Terminal Server

[Service]
ExecStart=/opt/butterfly/server.sh
